# host-agent/systemd/csm-agent.socket
[Unit]
Description=Claude Session Manager host agent socket

[Socket]
ListenStream=%t/csm-agent.sock
SocketMode=0660
SocketGroup=andres
Accept=yes

[Install]
WantedBy=sockets.target
